# etc/rc.local
#!/bin/bash

# Give ourselves some swap to deal with RAM issues
if [ -f /swapfile ]; then
  swapon /swapfile
fi
